# dev/systemd/avialeta.target
[Unit]
Description=Avialeta
Requires=avialeta-cache.service avialeta-api.service avialeta-ui.service

[Install]
WantedBy=default.target